# configs/boks.crontab
BOKS=/home/ubuntu/nt/boks/

@daily /usr/bin/node $BOKS/bin/dbbk.js
